# aircast/rootfs/etc/s6-overlay/s6-rc.d/esphome-bridge/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: AirCast
# Runs the ESPHome AirPlay Bridge via Shairport-Sync (if enabled)
# ==============================================================================

# Check if ESPHome support is enabled
if bashio::config.false 'esphome_enabled'; then
    bashio::log.info 'ESPHome support is disabled, skipping ESPHome bridge'
    exec sleep infinity
fi

bashio::log.info 'Starting ESPHome AirPlay Bridge (Shairport-Sync)...'

# Start Avahi daemon for mDNS (required for AirPlay discovery)
bashio::log.info 'Starting Avahi daemon...'
avahi-daemon --daemonize

# Run the Shairport-Sync manager for ESPHome devices
exec /usr/bin/python3 /usr/bin/shairport-manager.py
