# aircast/rootfs/etc/s6-overlay/s6-rc.d/init-aircast/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community App: AirCast
# Checks latency settings before starting the AirCast server
# ==============================================================================
declare latency

# Migrate app data from the Home Assistant config folder,
# to the app configuration folder.
if ! bashio::fs.file_exists '/config/aircast.xml' \
    && bashio::fs.file_exists '/homeassistant/aircast.xml'; then
    mv /homeassistant/aircast.xml /config/aircast.xml \
        || bashio::exit.nok "Failed to migrate Aircast configuration"
fi

# Create a configuration file, if it does not exist yet
if ! bashio::fs.file_exists '/config/aircast.xml'; then
    cp /etc/aircast.xml /config/aircast.xml
fi

# Warn if latency is below 500ms
latency=$(bashio::config 'latency_rtp')
if [[ "${latency}" -lt 500 && "${latency}" -ne 0 ]]; then
    bashio::log.warning \
        'Setting the RTP latency of AirPlay audio below 500ms is not recommended!'
fi
